# systemd/push-ddns.service
[Service]
WorkingDirectory=/opt/push-ddns
ExecStart=/opt/push-ddns/push-ddns.sh
User=ddns-pusher
